# deploy/tco-stopinfo-api.service
[Unit]
Description=TCO Stopinfo API (MQTT to HTTP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tco
Group=tco
WorkingDirectory=/opt/tco-stopinfo-api
Environment=TCO_CONFIG_FILE=/etc/tco-stopinfo-api/config.yaml
ExecStart=/opt/tco-stopinfo-api/.venv/bin/tco-stopinfo-api -c /etc/tco-stopinfo-api/config.yaml
Restart=on-failure
RestartSec=5
# Allow many concurrent connections from on-board displays
LimitNOFILE=65535

# Hardening (logs go to journald via stdout — no file log path needed)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
